# systemd/inky-image-display.service
[Unit]
Description=NASA Image of the Day Inky display
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=misterlamp
WorkingDirectory=/home/misterlamp
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/misterlamp/.config/inky-image-display.env
ExecStart=/home/misterlamp/.virtualenvs/pimoroni/bin/python /home/misterlamp/inky_image_display.py
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
